# conf/clusctld.service
[Unit]
Description=compute node control service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/clusctld
Restart=on-failure

[Install]
WantedBy=multi-user.target
